# supervisord.conf
[supervisord]
nodaemon=true

[program:record]
command=/bin/bash /app/record.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/record.err.log
stdout_logfile=/var/log/record.out.log